# .config/systemd/user/archive-trash.service
[Unit]
Description=Trash archiving service

[Service]
Type=oneshot
ExecStart=/bin/bash %h/.config/scripts/archive-files/archive-deleted-files.sh
