# .config/systemd/user/archive-trash.timer
[Unit]
Description=Trash archiving timer that runs once a day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
